# ops/systemd/hlauto.service
[Unit]
Description=Hyperliquid Auto Trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hlauto
Group=hlauto
WorkingDirectory=/opt/hlauto/trade
EnvironmentFile=/opt/hlauto/trade/.env.local
Environment=NODE_ENV=production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash /opt/hlauto/trade/ops/scripts/run-bot.sh
Restart=always
RestartSec=10
TimeoutStopSec=60
KillSignal=SIGINT
KillMode=process
UMask=0002
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/hlauto/trade/data

[Install]
WantedBy=multi-user.target
